# init.d.Linux/sysctl
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

depend()
{
	use hostname
	before bootmisc logger
	keywords novps
}

start()
{
	[ -e /etc/sysctl.conf ] || return 0
	ebegin "Configuring kernel parameters"
	sysctl -p /etc/sysctl.conf >/dev/null
	eend $?
}
