# init.d.Linux/volumes
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

VOLUME_ORDER=${VOLUME_ORDER:-${RC_VOLUME_ORDER:-raid evms lvm dm}}

# Dependancy information is in /etc/conf.d/volumes

start() {
	local x=

	# Start our volumes, RAID, LVM, etc
   	for x in ${VOLUME_ORDER} ; do
		start_addon "${x}"
   	done

	# Setup dm-crypt mappings if any
	start_addon dm-crypt

	return 0
}

stop() {
	local x= rev=

	# Stop dm-crypt mappings if any
	stop_addon dm-crypt
	stop_addon truecrypt

	# Stop our volumes, RAID, LVM, etc
   	for x in ${VOLUME_ORDER} ; do
		rev="${x} ${rev}"
	done

	for x in ${rev} ; do
		stop_addon "${x}"
   	done

	return 0
}

# vim: set ts=4 :
