# init.d.BSD/clock
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

opts="save"

description="Sets the local clock to UTC or Local Time." 
description_save="Saves the current time in the BIOS."

depend() {
	# BSD adjkerntz needs to be able to write to /etc
	if [ "${CLOCK}" = "UTC" -a -e /etc/wall_cmos_clock ] ||
	   [ "${CLOCK}" != "UTC" -a ! -e /etc/wall_cmos_clock ] ; then
		need checkroot
	fi
}

start() {
	if [ "${TIMEZONE-Factory}" = "Factory" ] ; then
		ewarn "Your TIMEZONE in /etc/conf.d/clock is still set to Factory!"
	fi

	local TBLURB="Local Time"
	[ "${CLOCK}" = "UTC" ] && TBLURB="UTC"
	ebegin "Starting the System Clock Adjuster [${TBLURB}]"
	if [ "${CLOCK}" != "UTC" ] ; then
		echo >/etc/wall_cmos_clock
		start-stop-daemon --start --exec /sbin/adjkerntz -- -i
	else
		rm -f /etc/wall_cmos_clock
		/sbin/adjkerntz -i
	fi
	eend $?
}

save() {
	local TBLURB="Local Time"
	[ "${CLOCK}" = "UTC" ] && TBLURB="UTC"
	ebegin "Setting hardware clock using the system clock" "[${TBLURB}]"
	adjkerntz -a
	eend $?
}

stop() {
    # Don't tweak the hardware clock on LiveCD halt.
	[ -z "${CDBOOT}" -a "${CLOCK_SYSTOHC}" = "yes" ] && save 

	ebegin "Stopping the System Clock Adjuster"
	if start-stop-daemon --test --quiet --stop --exec /sbin/adjkerntz ; then
		start-stop-daemon --stop --exec /sbin/adjkerntz
		eend $?
	else
		eend 0
	fi
}

# vim: set ts=4 :
