# init.d.BSD/moused
#!/sbin/runscript
# Copyright 2007 Roy Marples
# All rights reserved

# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions
# are met:
# 1. Redistributions of source code must retain the above copyright
#    notice, this list of conditions and the following disclaimer.
# 2. Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE AUTHOR AND CONTRIBUTORS ``AS IS'' AND
# ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED.  IN NO EVENT SHALL THE AUTHOR OR CONTRIBUTORS BE LIABLE
# FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL
# DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS
# OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION)
# HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT
# LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY
# OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF
# SUCH DAMAGE.

mouse=${SVCNAME##*.}
if [ -n "${name}" -a "${mouse}" != "moused" ] ; then
	moused_device=/dev/"${mouse}"
	pidfile=/var/run/moused-"${mouse}".pid
else
	pidfile=/var/run/moused.pid
fi
name="Console Mouse Daemon"
[ -n "${moused_device}" ] && name="${name} (${moused_device})"

depend() {
	need localmount
	after bootmisc
}

start() {
	ebegin "Starting ${name}" 

	if [ -z "${moused_device}" ] ; then
		local dev=
		for dev in /dev/psm[0-9]* /dev/ums[0-9]* ; do
			[ -c "${dev}" ] || continue
			[ -e /var/run/moused-"${dev##*/}".pid ] && continue 
			moused_device=${dev}
			eindent
			einfo "Using mouse on ${moused_device}"
			eoutdent
			break
		done
	fi

	if [ -z "${moused_device}" ] ; then
		eend 1 "No mouse device found" 
		return 1
	fi

	start-stop-daemon --start --exec /usr/sbin/moused \
		--pidfile "${pidfile}" \
		-- ${moused_args} -p "${moused_device}" -I "${pidfile}"
	local retval=$?
	
	local ttyv=
	for ttyv in /dev/ttyv*; do
		vidcontrol < "${ttyv}" -m on
		: $((retval+= $?))
	done
	
	eend ${retval} "Failed to start moused"
}

# vim: set ts=4 :
