# init.d.BSD/newsyslog
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

required_files="/etc/newsyslog.conf"

depend()
{
	need localmount
}

start()
{
	ebegin "Creating and/or trimming log files"
	newsyslog ${newsyslog_args}
	eend $?
}
