# init.d.BSD/rpcbind
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=/usr/sbin/rpcbind
command_args=${rpcbind_args}
name="RPC program number mapper"

depend()
{
	provide rpc
	need localmount
	use net logger dns
	before inetd xinetd ntpd ntp-client
}

stop_post()
{
	# rpcbind returns too fast, so sleep for a second
	sleep 1
}
