# init.d.BSD/sysctl.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

depend()
{
	use hostname
	before bootmisc logger
	keyword noprefix
}

start()
{
	[ -e /etc/sysctl.conf ] || return 0
	
	ebegin "Configuring kernel parameters"
	local retval=0 var= comments=
	while read var comments; do
		case "${var}" in
			""|"#"*) continue;;
		esac
		sysctl -w "${var}" >/dev/null || retval=1
	done < /etc/sysctl.conf
	eend ${retval} "Some errors were encountered"
	eend $?
}
