# init.d.BSD/syslogd
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=/usr/sbin/syslogd
command_args=${syslogd_args}
pidfile=/var/run/syslog.pid
name="System Logger Daemon"

depend() {
	provide logger
	use net
	need localmount
	after bootmisc newsyslog
}
