# init.d.BSD/syslogd
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

command=/usr/sbin/syslogd
command_args=${syslogd_args}
pidfile=$(strings ${command} | grep /var/run/syslog)
name="System Logger Daemon"

depend()
{
	provide logger
	use net
	need localmount
	after bootmisc newsyslog
}
