# init.d.misc/avahid.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=@PREFIX@%%/sbin/avahi-daemon
command_args="${avahid_args} -D"
pidfile=@VARBASE@/run/avahi-daemon/pid
name="Avahi Service Advertisement Daemon"

depend() {
	use dns
	need localmount dbus
	after bootmisc
}
