# init.d.misc/avahi-dnsconfd.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=@PREFIX@/sbin/avahi-dnsconfd
command_args="${avahi_dnsconfd_args} -D"
pidfile=@VARBASE@/run/avahi-dnsconfd.pid
name="Avahi DNS Configuration Daemon"

depend() {
	use dns
	need localmount dbus
	after bootmisc
}
