# init.d.misc/dnsmasq.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=@PREFIX@/sbin/dnsmasq
command_args=${dnsmasq_args}
pidfile=@VARBASE@/run/dnsmasq.pid
required_files=/etc/dnsmasq.conf

extra_started_commands="reload"

depend() {
	provide dns
	need localmount net
	after bootmisc
}

reload() {
	ebegin "Reloading ${SVCNAME}"
	start-stop-daemon --signal SIGHUP --pidfile "${pidfile}"
	eend $?
}

# vim: set ts=4 :
