# init.d.misc/polkitd.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples
# All rights reserved

command=@PREFIX@/sbin/polkitd
pidfile=@VARBASE@/run/polkitd/polkitd.pid
command_args="${polkitd_args}"
name="PolicyKit Daemon"

depend() {
	need dbus
}

# vim: set ts=4 :
