# etc/profile.d/Makefile
DIR   = /etc/profile.d

TOPDIR = ../..
include $(TOPDIR)/default.mk
