# init.d/hostname.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

description="Sets the hostname of the machine."

depend() {
	keyword noprefix
}

start()
{
	hostname=${hostname-${HOSTNAME-localhost}}
	ebegin "Setting hostname to ${hostname}"
	hostname "${hostname}"
	eend $? "Failed to set the hostname"
}
