# init.d/local.in
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

description="Executes user commands in /etc/conf.d/local"

depend()
{
	after *
	keyword notimeout
}

start()
{
	ebegin "Starting local"

	if type local_start >/dev/null 2>&1; then
		local_start
	fi

	# Support old configs
	if [ -e /etc/conf.d/local.start ]; then
		. /etc/conf.d/local.start
	fi

	eend $? "Failed to start local"
}

stop()
{
	ebegin "Stopping local"

	if type local_start >/dev/null 2>&1; then
		local_stop
	fi

	# Support old configs
	if [ -e /etc/conf.d/local.stop ]; then
		. /etc/conf.d/local.stop
	fi

	eend $? $"Failed to stop local"
}
