# init.d/swap
#!/sbin/runscript
# Copyright 2007-2008 Roy Marples <roy@marples.name>
# All rights reserved. Released under the 2-clause BSD license.

depend()
{
	need localmount
	keyword nojail
}

start()
{
	ebegin "Activating swap"
	swapon -a >/dev/null
	eend 0 # If swapon has nothing todo it errors, so always return 0
}	

stop()
{
	ebegin "Deactivating swap"
	swapoff -a >/dev/null
	eend 0 
}
