# init.d/urandom
#!/sbin/runscript
# Copyright 2007 Roy Marples
# All rights reserved

# Redistribution and use in source and binary forms, with or without
# modification, are permitted provided that the following conditions
# are met:
# 1. Redistributions of source code must retain the above copyright
#    notice, this list of conditions and the following disclaimer.
# 2. Redistributions in binary form must reproduce the above copyright
#    notice, this list of conditions and the following disclaimer in the
#    documentation and/or other materials provided with the distribution.
#
# THIS SOFTWARE IS PROVIDED BY THE AUTHOR AND CONTRIBUTORS ``AS IS'' AND
# ANY EXPRESS OR IMPLIED WARRANTIES, INCLUDING, BUT NOT LIMITED TO, THE
# IMPLIED WARRANTIES OF MERCHANTABILITY AND FITNESS FOR A PARTICULAR PURPOSE
# ARE DISCLAIMED.  IN NO EVENT SHALL THE AUTHOR OR CONTRIBUTORS BE LIABLE
# FOR ANY DIRECT, INDIRECT, INCIDENTAL, SPECIAL, EXEMPLARY, OR CONSEQUENTIAL
# DAMAGES (INCLUDING, BUT NOT LIMITED TO, PROCUREMENT OF SUBSTITUTE GOODS
# OR SERVICES; LOSS OF USE, DATA, OR PROFITS; OR BUSINESS INTERRUPTION)
# HOWEVER CAUSED AND ON ANY THEORY OF LIABILITY, WHETHER IN CONTRACT, STRICT
# LIABILITY, OR TORT (INCLUDING NEGLIGENCE OR OTHERWISE) ARISING IN ANY WAY
# OUT OF THE USE OF THIS SOFTWARE, EVEN IF ADVISED OF THE POSSIBILITY OF
# SUCH DAMAGE.

urandom_seed=${urandom_seed:-${URANDOM_SEED:-/var/run/random-seed}}

description="Initializes the random number generator."

depend() {
	need localmount
}

save_seed() {
	local poolsize=1

	if [ -e /proc/sys/kernel/random/poolsize ]; then
		poolsize=$(($(cat /proc/sys/kernel/random/poolsize) / 4096))
	fi

	(	# sub shell to prevent umask pollution
		umask 077
		dd if=/dev/urandom of="${urandom_seed}" count=${poolsize} 2>/dev/null
	)
}

start() {
	[ -c /dev/urandom ] || return
	if [ -f "${urandom_seed}" ]; then
		ebegin "Initializing random number generator"
		cat "${urandom_seed}" > /dev/urandom
		eend $? "Error initializing random number generator"
	fi
	rm -f "${urandom_seed}" && save_seed
	return 0
}

stop() {
	ebegin "Saving random seed"
	save_seed
	eend $? "Failed to save random seed"
}

# vim: set ts=4 :
